# systemd/bluetooth.service.d/nobattery.conf
# Disable battery plugin for PiAware BLE configuration workaround

[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --noplugin=battery
